# dom0-updates/qubes-dom0-updates.cron
#!/bin/bash

qvm-features dom0 last-updates-check "$(date +'%Y-%m-%d %T')"

if [ "$(qvm-features dom0 service.qubes-update-check || echo 1)" != 1 ]; then
    exit 0
fi

# If no updates available - exit here
qubes-dom0-update --check-only >/dev/null && exit
RETCODE=$?

if [ "$RETCODE" -ne 100 ]; then
    echo "ERROR: Error checking for updates" >&2
    exit $RETCODE
fi

qvm-features dom0 updates-available 1
